# packaging/DEBIAN/postinst
#!/bin/bash
set -e

# Create system user and group
if ! getent passwd zoneweaver > /dev/null 2>&1; then
    adduser --system --group --home /opt/zoneweaver \
            --no-create-home --shell /usr/sbin/nologin \
            --gecos "Zoneweaver system user" zoneweaver
fi

# Set ownership and permissions
chown -R zoneweaver:zoneweaver /opt/zoneweaver
chown -R zoneweaver:zoneweaver /var/lib/zoneweaver
chown -R zoneweaver:zoneweaver /var/log/zoneweaver
chown -R zoneweaver:zoneweaver /etc/zoneweaver


# Generate JWT secret if it doesn't exist
if [ ! -f /etc/zoneweaver/.jwt-secret ]; then
    openssl rand -hex 32 > /etc/zoneweaver/.jwt-secret
    chown zoneweaver:zoneweaver /etc/zoneweaver/.jwt-secret
    chmod 600 /etc/zoneweaver/.jwt-secret
    echo "Generated JWT secret in /etc/zoneweaver/.jwt-secret"
fi

# Update config file with actual JWT secret
if [ -f /etc/zoneweaver/.jwt-secret ]; then
    JWT_SECRET=$(cat /etc/zoneweaver/.jwt-secret)
    sed -i "s/__JWT_SECRET_FROM_FILE__/${JWT_SECRET}/g" /etc/zoneweaver/config.yaml
    echo "Updated config with JWT secret"
fi

# Initialize database directory
mkdir -p /var/lib/zoneweaver/database
chown zoneweaver:zoneweaver /var/lib/zoneweaver/database

# Create SSL directory for certificate generation
mkdir -p /etc/zoneweaver/ssl
chown zoneweaver:zoneweaver /etc/zoneweaver/ssl
chmod 700 /etc/zoneweaver/ssl

# Reload systemd daemon
systemctl daemon-reload

echo ""
echo "======================================================================"
echo "Zoneweaver has been installed successfully!"
echo ""
echo "Next steps:"
echo "1. Review configuration: /etc/zoneweaver/config.yaml"
echo "2. Start the service: systemctl enable --now zoneweaver"
echo "3. Check status: systemctl status zoneweaver"
echo "4. View logs: journalctl -u zoneweaver -f"
echo ""
echo "Default access: https://localhost:3443 (with self-signed certificate)"
echo "======================================================================"
echo ""

#DEBHELPER#
