# packaging/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable service before removal
if systemctl is-active --quiet zoneweaver; then
    echo "Stopping Zoneweaver service..."
    systemctl stop zoneweaver
fi

if systemctl is-enabled --quiet zoneweaver; then
    echo "Disabling Zoneweaver service..."
    systemctl disable zoneweaver
fi

#DEBHELPER#
